# ndncert/debian/ndncert-server.upstart
description "NDN Certificate Management System Library"
author      "Zhiyi Zhang <zhiyi@cs.ucla.edu>"

start on started  nfd
stop  on stopping nfd

respawn
respawn limit unlimited

setuid ndn
setgid ndn

script
  FLAGS = ""
  if [ -f /etc/default/ndncert ]; then
      . /etc/default/ndncert
  fi

  export HOME=/var/lib/ndn/ndncert
  exec /usr/bin/ndncert-ca-server -c /etc/ndncert/ndncert-ca.conf $FLAGS
end script
