# nfd/debian/ndn-autoconfig-client.upstart
description "NDN Autoconfig Server"
author      "Alexander Afanasyev <alexander.afanasyev@ucla.edu>"

start on started  nfd
stop  on stopping nfd

respawn
respawn limit unlimited

setuid ndn
setgid ndn

script
  if [ -f /etc/ndn/autoconfig.conf ]; then
    exec /usr/bin/ndn-autoconfig -d -c /etc/ndn/autoconfig.conf
  fi
end script
