# nfd/debian/nfd-status-http-server.upstart
description "NDN HTTP Status Server"
author      "Alexander Afanasyev <alexander.afanasyev@ucla.edu>"

start on started  nfd
stop  on stopping nfd

respawn
respawn limit unlimited

script
  FLAGS=""
  if [ -f /etc/default/nfd-status-http-server ]; then
      . /etc/default/nfd-status-http-server
  fi
  export HOME=/tmp/nfd-status-http-server
  exec /usr/bin/nfd-status-http-server $FLAGS
end script
