# nfd/debian/nfd.nfd-restart.upstart
# nfd-restart.conf
#
# Restarting NDN Forwarding Daemon on network connectivity changes

start on net-device-up
task

script
        status nfd | grep -q start/ || stop
        restart nfd
end script
